# overlay/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/bin/sh -c "cd /srv/wyper/src && bash main.sh"
StandardInput=tty
StandardOutput=tty

