# overlay/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/bin/sh -c "less +F /tmp/wyper/alog || sleep 2s"
StandardInput=tty
StandardOutput=tty

