# overlay/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/bin/sh -c "sleep 10s && tail -fn 1000 /tmp/wyper/alog"
StandardInput=tty
StandardOutput=tty

